# deploy/chatbot.service
[Unit]
Description=Stadfirma Chatbot FastAPI Application
After=network.target

[Service]
Type=simple
User=chatbot
Group=chatbot
WorkingDirectory=/opt/chatbot/app
Environment=PATH=/opt/chatbot/app/venv/bin
ExecStart=/opt/chatbot/app/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
